# .config/systemd/user/sommelier-x@0.service.d/override.conf
[Service]
Environment="SOMMELIER_FRAME_COLOR=#3c3f41"
